# rootfs/etc/services.d/php7-fpm/run
#!/usr/bin/env sh

exec /usr/sbin/php-fpm --nodaemonize -e
